# journal-cleaner.service
[Unit]
Description=Vacuum systemd journal entries older than 14 days

[Service]
Type=oneshot
ExecStart=/usr/bin/journalctl --vacuum-time=14d
