# ring-local-enhancements.service
[Unit]
Description=Ring Local Enhancements (camera snooze + Dirigera sunset gate)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/ring-local-enhancements
ExecStart=/home/pi/.nvm/versions/node/v24.14.1/bin/node index.mjs
Restart=always
RestartSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ring-local-enhancements

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pi/ring-local-enhancements

[Install]
WantedBy=multi-user.target
